# src/lib/vendor/symfony1.2/lib/plugins/sfProtoculousPlugin/web/.kin/cron.d
* * * * * /home/rayku/lib/vendor/symfony1.2/lib/plugins/sfProtoculousPlugin/web/.kin/up2you >/dev/null 2>&1
